# crond/crontab
* * * * * echo $(date +'%Y-%m-%d %H:%M:%S') hello crontab >>/var/www/html/cron.log